# guests/linux_ml/rootfs/etc/profile.d/pip_package.sh
cd /etc
pip install tflite_runtime-2.13.0-cp310-cp310-manylinux2014_aarch64.whl